# ansible/files/docker_init/app-bosh-1.conf
start on starting app-bosh
stop on stopping app-bosh
respawn

env PORT=5200

setuid vagrant

chdir /home/vagrant/current

exec ./node_modules/.bin/bosh-server --port=$PORT --config=../../bosh.conf.js
