# ansible/files/docker_init/app-chat-1.conf
start on starting app-chat
stop on stopping app-chat
respawn

env PORT=5100

setuid vagrant

chdir /home/vagrant/current

exec node chat.js
