# ansible/files/docker_init/attalos-chat.conf
start on attalos-chat
stop on attalos-chat

respawn

env PORT=5100

setuid vagrant

chdir /home/vagrant/current

exec node chat.js
